# water-on.timer
[Unit]
Description=Run water-on daily at 20:00

[Timer]
OnCalendar=*-*-* 20:00:00
Persistent=true

[Install]
WantedBy=timers.target
